# systemd/owncloud-uploader.service
[Unit]
Description=owncloud uploadr

[Service]
ExecStart=__HOME__/bin/owncloud-uploader.pl
Restart=on-failure

[Install]
WantedBy=default.target
